# rpi/rc.local
# append this to the end of /etc/rc.local

python /home/rpi/startup_blink.py &

python /home/rpi/bulbs.py &
python /home/rpi/rotary_tracks.py &
python /home/rpi/toggle_restart.py &
python /home/rpi/toggle_shutdown.py &
